# deploy/daw-ai.service
[Unit]
Description=DAW-AI local prompt-driven music studio
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/lib/daw-ai
ExecStart=/usr/local/bin/daw-ai --port 8888
Restart=on-failure
RestartSec=5s
Environment=PATH=/usr/local/bin:/usr/bin:/bin
LoadCredential=gemini-api-key:/etc/daw-ai/gemini-api-key
StateDirectory=daw-ai
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
